# package/siflower/bin/led-button/files/led_button_init
#!/bin/sh /etc/rc.common
# Copyright (c) 2011-2012 OpenWrt.org

START=50

USE_PROCD=1
PROG=/bin/led-button


set_led_button() {
	procd_open_instance
	procd_set_param command "$PROG"
	procd_append_param command "-b" "$1"
	procd_set_param respawn
	procd_close_instance
}

start_service() {
	if [ -d "/dev/input/" ];then
		for i in `ls /dev/input/`
		do
			input="/dev/input/$i"
			set_led_button $input
		done
	else
		echo "error >>> no input event node !!!" >  /dev/console
	fi

}
