# package/siflower/bin/utils/files/base-files/etc/hotplug.d/net/01-wlanctl
#!/bin/sh

if [ "$ACTION" = add ]; then
	[ "x$DEVICENAME" = "xwlan0-guest" ] && gwifi speed_rst 0
	[ "x$DEVICENAME" = "xwlan1-guest" ] && gwifi speed_rst 1
fi
