# package/siflower/luci-siflower/themes/luci-theme-siwifi/root/etc/uci-defaults/luci-theme-siwifi
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Siwifi=/luci-static/siwifi.org
	set luci.main.mediaurlbase=/luci-static/siwifi.org
        commit luci
EOF

