# package/mail/fetchmail/files/fetchmail.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2008-2023 OpenWrt.org

START=99

start() {
        service_start /usr/bin/fetchmail_set
}

