# package/network/services/luci-app-i2c-ssd1306/root/etc/uci-defaults/i2c-ssd1306
#!/bin/sh

/etc/init.d/i2c-ssd1306 enable
/etc/init.d/i2c-ssd1306 start

rm -f /tmp/luci-indexcache
exit 0
