# package/network/services/luci-app-i2c-lcd1602/root/etc/uci-defaults/i2c-lcd1602
#!/bin/sh

/etc/init.d/i2c-lcd1602 enable
/etc/init.d/i2c-lcd1602 start

rm -f /tmp/luci-indexcache
exit 0
